# configs/etc/supervisor/conf.d/web_test.conf
[supervisord]
nodaemon=true

[program:web_test]
command=python3 hello.py
autostart=true
autorestart=true
startsecs=2
user=root
